# plasma/dot-config/autostart/conky.desktop
[Desktop Entry]
Exec=conky -qdc .config/conky/way-clock.conf
Icon=conky-logomark-violet
Name=conky
Path=
Terminal=false
Type=Application
